# 95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm",  RUN+="/usr/bin/systemctl --no-block start fix-monitor-font-scale.service"

