# fix-monitor-font-scale.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix_font_scale.sh
User=borovaka
Group=borovaka
